# root/etc/s6-overlay/s6-rc.d/service-forwarder/run
#!/command/with-contenv bash
# shellcheck shell=bash

# shellcheck source=/dev/null
source /etc/s6-overlay/scripts/bash-functions

umask "${UMASK}"

if [[ -n ${VPN_AUTO_PORT_FORWARD_TO_PORTS} ]]; then
    VPN_PORT_REDIRECTS=${VPN_AUTO_PORT_FORWARD_TO_PORTS}
    log_wrn "[FORWARDER] The variable [VPN_AUTO_PORT_FORWARD_TO_PORTS] is deprecated in favor of [VPN_PORT_REDIRECTS], please adjust your config!"
fi

# Wait for other services to do their thing
sleep 10

old_port="0"
while true; do
    re='^[0-9]+$'
    if [[ "${VPN_AUTO_PORT_FORWARD}" =~ $re ]]; then
        echo "${VPN_AUTO_PORT_FORWARD}" > "${CONFIG_DIR}/wireguard/forwarded_port"
        sleep_time=86400
    else
        sleep_time=60
    fi
    if [[ -f "${CONFIG_DIR}/wireguard/forwarded_port" ]]; then
        port=$(sed -n '1p' "${CONFIG_DIR}/wireguard/forwarded_port")
        [[ "${port}" =~ $re ]] || \
            log_err "[FORWARDER] The file [${CONFIG_DIR}/wireguard/forwarded_port] contains a value that's not a number!"
        if [[ "${old_port}" != "${port}" ]] && [[ "${port}" =~ $re ]]; then
            log_inf "[FORWARDER] Forwarded port is [${port}]."
            IFS=',' read -ra forward_ports <<< "${VPN_PORT_REDIRECTS%,}"
            for forward_port in "${forward_ports[@]}"; do
                [[ ${forward_port} == *@* ]] && continue
                port_is_redirected="true"
                forward_port="${forward_port//[[:space:]]/}"
                internal="${forward_port#*@}"
                del_rule_redirect_port "${old_port}" "${internal}"
                add_rule_redirect_port "${port}" "${internal}"
                add_rule_expose_port "${internal}"
            done
            if [[ "${port_is_redirected}" != "true" ]]; then
                del_rule_expose_port "${old_port}/tcp"
                del_rule_expose_port "${old_port}/udp"
                add_rule_expose_port "${port}/tcp"
                add_rule_expose_port "${port}/udp"
            fi
            # shellcheck source=/dev/null
            source "${APP_DIR}/forwarded_port_update_app"
            if [[ -f "${CONFIG_DIR}/wireguard/${VPN_CONF}-port.sh" ]]; then
                log_inf "[FORWARDER] Script [${CONFIG_DIR}/wireguard/${VPN_CONF}-port.sh] found. Executing..."
                bash "${CONFIG_DIR}/wireguard/${VPN_CONF}-port.sh"
            fi
            old_port="${port}"
        fi
    else
        log_wrn "[FORWARDER] The file [${CONFIG_DIR}/wireguard/forwarded_port] was not found."
    fi
    sleep "${sleep_time}"
done
