# root/etc/s6-overlay/s6-rc.d/service-forwarder/run
#!/command/with-contenv bash
# shellcheck shell=bash

# shellcheck source=/dev/null
source /etc/s6-overlay/scripts/general-functions

umask "${UMASK}"

add_rule_expose_port() {
    internal="${1}"
    new_rule=$(grep "iifname \"${VPN_CONF}\" ${internal##*/} dport ${internal%/*} " <<< "$(nft --handle list chain inet hotio input)" | sed 's/.* # handle //g')
    if [[ -z "${new_rule}" ]]; then
        nft add rule inet hotio input iifname "${VPN_CONF}" "${internal##*/}" dport "${internal%/*}" counter accept
        log_inf "[VPN] Exposed port [${internal%/*}/${internal##*/}] on interface [${VPN_CONF}]."
    fi
}

del_rule_expose_port() {
    internal="${1}"
    old_rule=$(grep "iifname \"${VPN_CONF}\" ${internal##*/} dport ${internal%/*} " <<< "$(nft --handle list chain inet hotio input)" | sed 's/.* # handle //g')
    if [[ -n "${old_rule}" ]]; then
        nft delete rule inet hotio input handle "${old_rule}"
        log_inf "[VPN] Deleted exposed port [${internal%/*}/${internal##*/}] on interface [${VPN_CONF}]."
    fi
}

add_rule_redirect_port() {
    forward_port="${1}"
    internal="${2}"
    new_rule=$(grep "iifname \"${VPN_CONF}\" ${internal##*/} dport ${forward_port} redirect to :${internal##*/} " <<< "$(nft --handle list chain inet hotio prerouting)" | sed 's/.* # handle //g')
    if [[ -z "${new_rule}" ]]; then
        nft add rule inet hotio prerouting iifname "${VPN_CONF}" "${internal##*/}" dport "${forward_port}" redirect to :"${internal%/*}"
        log_inf "[VPN] Redirected port [${forward_port}/${internal##*/}] to [${internal%/*}/${internal##*/}] on interface [${VPN_CONF}]."
    fi
}

del_rule_redirect_port() {
    forward_port="${1}"
    internal="${2}"
    old_rule=$(grep "iifname \"${VPN_CONF}\" ${internal##*/} dport ${forward_port} redirect to :${internal%/*} " <<< "$(nft --handle list chain inet hotio prerouting)" | sed 's/.* # handle //g')
    if [[ -n "${old_rule}" ]]; then
        nft delete rule inet hotio prerouting handle "${old_rule}"
        log_inf "[VPN] Deleted redirected port [${forward_port}/${internal##*/}] to [${internal%/*}/${internal##*/}] on interface [${VPN_CONF}]."
    fi
}

if [[ -n ${VPN_AUTO_PORT_FORWARD_TO_PORTS} ]]; then
    VPN_PORT_REDIRECTS=${VPN_AUTO_PORT_FORWARD_TO_PORTS}
    log_wrn "[VPN] The variable [VPN_AUTO_PORT_FORWARD_TO_PORTS] is deprecated in favor of [VPN_PORT_REDIRECTS], please adjust your config!"
fi

IFS=',' read -ra forward_ports <<< "${VPN_PORT_REDIRECTS%,}"
for forward_port in "${forward_ports[@]}"; do
    forward_port="${forward_port// /}"
    internal="${forward_port##*@}"
    add_rule_expose_port "${internal}"
    if [[ ${forward_port} =~ "@" ]] && [[ "${forward_port%@*}" != "${internal%/*}" ]]; then
        add_rule_redirect_port "${forward_port%@*}" "${internal}"
    fi
done

re='^[0-9]+$'
if [[ "${VPN_AUTO_PORT_FORWARD}" =~ $re ]] ; then
    echo "${VPN_AUTO_PORT_FORWARD}" > "${CONFIG_DIR}/wireguard/forwarded_port"
fi

sleep 10

old_port="0"
while true; do
    if [[ -f "${CONFIG_DIR}/wireguard/forwarded_port" ]] && [[ "${VPN_AUTO_PORT_FORWARD}" != "false" ]]; then
        port=$(sed -n '1p' "${CONFIG_DIR}/wireguard/forwarded_port")
        if [[ "${old_port}" != "${port}" ]]; then
            log_inf "[VPN] Forwarded port is [${port}]."
            IFS=',' read -ra forward_ports <<< "${VPN_PORT_REDIRECTS%,}"
            for forward_port in "${forward_ports[@]}"; do
                forward_port="${forward_port// /}"
                internal="${forward_port##*@}"
                if [[ ! ${forward_port} =~ "@" ]]; then
                    port_is_redirected="true"
                    del_rule_redirect_port "${old_port}" "${internal}"
                    [[ "${port}" != "${internal%/*}" ]] && \
                        add_rule_redirect_port "${port}" "${internal}"
                fi
            done
            if [[ "${port_is_redirected}" != "true" ]]; then
                del_rule_expose_port "${old_port}/tcp"
                del_rule_expose_port "${old_port}/udp"
                add_rule_expose_port "${port}/tcp"
                add_rule_expose_port "${port}/udp"
            fi
            # shellcheck source=/dev/null
            source "${APP_DIR}/forwarded_port_update_app"
            if [[ -f "${CONFIG_DIR}/wireguard/${VPN_CONF}-port.sh" ]]; then
                log_inf "[VPN] Script [${CONFIG_DIR}/wireguard/${VPN_CONF}-port.sh] found. Executing..."
                bash "${CONFIG_DIR}/wireguard/${VPN_CONF}-port.sh"
            fi
            old_port="${port}"
        fi
    fi
    sleep 60
done
