# root/etc/s6-overlay/s6-rc.d/service-pia/run
#!/command/with-contenv bash
# shellcheck shell=bash

# shellcheck source=/dev/null
source /etc/s6-overlay/scripts/pia-functions

umask "${UMASK}"

while true; do
    if [[ -f "${CONFIG_DIR}/wireguard/${VPN_CONF}.conf" ]]; then
        pia_server_hostname=$(grep '# hostname: ' < "${CONFIG_DIR}/wireguard/${VPN_CONF}.conf" | awk -F ': ' '{print $2}')
        pia_server_ip=$(grep '# ip: ' < "${CONFIG_DIR}/wireguard/${VPN_CONF}.conf" | awk -F ': ' '{print $2}')
        if [[ "${VPN_PIA_PORT_FORWARD_PERSIST}" == "true" ]] && [[ -f "${CONFIG_DIR}/wireguard/${VPN_CONF}.persist.json" ]]; then
            echo "[INF] [$(date '+%Y-%m-%d %H:%M:%S')] [PIA] Using persistent forwarded port."
            pia_response=$(cat "${CONFIG_DIR}/wireguard/${VPN_CONF}.persist.json")
        else
            echo "[INF] [$(date '+%Y-%m-%d %H:%M:%S')] [PIA] Fetching new forwarded port..."
            pia_response=$(curl --retry 5 --retry-max-time 60 --max-time 10 --interface "${VPN_CONF}" -Gs --connect-to "$pia_server_hostname::$pia_server_ip:" --cacert "${CONFIG_DIR}/wireguard/pia.ca.rsa.4096.crt" --data-urlencode "token=$(get_token)" "https://${pia_server_hostname}:19999/getSignature")
        fi
        if [[ $(echo "${pia_response}" | jq -re '.status') == "OK" ]]; then
            [[ "${VPN_PIA_PORT_FORWARD_PERSIST}" == "true" ]] && echo "${pia_response}" > "${CONFIG_DIR}/wireguard/${VPN_CONF}.persist.json"
            pia_payload=$(jq -re '.payload' <<< "${pia_response}")
            pia_signature=$(jq -re '.signature' <<< "${pia_response}")
            pia_payload_decoded=$(base64 -d <<< "${pia_payload}")
            pia_port=$(jq -re '.port' <<< "${pia_payload_decoded}")
            pia_port_expiration=$(jq -re '.expires_at' <<< "${pia_payload_decoded}")
            [[ $(date -d "${pia_port_expiration}" +%s) -lt $(date -d "$(date)" +%s) ]] && rm -rf "${CONFIG_DIR}/wireguard/${VPN_CONF}.persist.json" && continue
            echo "[INF] [$(date '+%Y-%m-%d %H:%M:%S')] [PIA] Forwarded port will expire at [${pia_port_expiration}]."
            echo "${pia_port}" > "${CONFIG_DIR}/wireguard/forwarded_port"
            while [[ $(date -d "${pia_port_expiration}" +%s) -gt $(date -d "$(date)" +%s) ]]; do
                bind_port_response="$(curl --retry 5 --retry-max-time 60 --max-time 10 --interface "${VPN_CONF}" -Gs --connect-to "$pia_server_hostname::$pia_server_ip:" --cacert "${CONFIG_DIR}/wireguard/pia.ca.rsa.4096.crt" --data-urlencode "payload=${pia_payload}" --data-urlencode "signature=${pia_signature}" "https://${pia_server_hostname}:19999/bindPort")"
                if ! [[ $(jq -re '.status' <<< "${bind_port_response}") == "OK" ]]; then
                    echo "[WRN] [$(date '+%Y-%m-%d %H:%M:%S')] [PIA] Port bind failed!"
                    break
                fi
                sleep 900
            done
        else
            echo "[WRN] [$(date '+%Y-%m-%d %H:%M:%S')] [PIA] Signature and Payload retrieval failed!"
        fi
        sleep 60
    else
        echo "[ERR] [$(date '+%Y-%m-%d %H:%M:%S')] [PIA] File [${CONFIG_DIR}/wireguard/${VPN_CONF}.conf] not found! Checking again in 60 seconds..."
        sleep 60
    fi
done
