# root/etc/s6-overlay/s6-rc.d/init-setup/run
#!/command/with-contenv bash
# shellcheck shell=bash

umask "${UMASK}"

mask() {
    local n=3
    [[ ${#1} -le 5 ]] && n=$(( ${#1} - 3 ))
    local a="${1:0:${#1}-n}"
    local b="${1:${#1}-n}"
    printf "%s%s\n" "${a//?/*}" "$b"
}

echo "
$(figlet "hotio")
Documentation: https://hotio.dev/$([[ -f /image_stats.json ]] && jq -r '.app' < /image_stats.json)
Support:       https://hotio.dev/discord
Image:         $([[ -f /image_stats.json ]] && jq -r '.image' < /image_stats.json)
Revision:      $([[ -f /image_stats.json ]] && jq -r '.revision' < /image_stats.json)
Version:       $([[ -f /image_stats.json ]] && jq -r '.version' < /image_stats.json)

----------------------------------------------------------------------
ENVIRONMENT BASE
----------------------------------------------------------------------
PUID=${PUID}
PGID=${PGID}
UMASK=${UMASK}
TZ=${TZ}
VPN_ENABLED=${VPN_ENABLED}
VPN_PROVIDER=${VPN_PROVIDER}
VPN_LAN_NETWORK=${VPN_LAN_NETWORK}
VPN_CONF=${VPN_CONF}
VPN_ADDITIONAL_PORTS=${VPN_ADDITIONAL_PORTS}
VPN_AUTO_PORT_FORWARD=${VPN_AUTO_PORT_FORWARD}
VPN_PIA_USER=$(mask "${VPN_PIA_USER}")
VPN_PIA_PASS=$(mask "${VPN_PIA_PASS}")
VPN_PIA_PREFERRED_REGION=${VPN_PIA_PREFERRED_REGION}
PRIVOXY_ENABLED=${PRIVOXY_ENABLED}
----------------------------------------------------------------------
"

echo "Executing usermod..."
mkdir "/tmp/temphome"
usermod -d "/tmp/temphome" hotio
usermod -o -u "${PUID}" hotio
usermod -d "${CONFIG_DIR}" hotio
rm -rf "/tmp/temphome"
groupmod -o -g "${PGID}" hotio

echo "Applying permissions to ${CONFIG_DIR}"
chmod "=rwx" "${CONFIG_DIR}"
chown hotio:hotio "${CONFIG_DIR}"
